# srv/salt/projects/templates/supervisor-uwsgi.conf
[program:{{ program_name }}]
autorestart=true
command={{ uwsgi_bin }} --ini {{ uwsgi_ini }}
redirect_stderr=true
